# testdata/supervisord.d/prometheus.ini
; Managed by pmm-managed. DO NOT EDIT.

[program:prometheus]
command = /bin/echo Prometheus is substituted by VictoriaMetrics
user = pmm
autorestart = false
autostart = false
